# files/etc/init.d/myvpn
#!/bin/sh /etc/rc.common
# /etc/init.d/myvpn
START=95
PID=/tmp/myvpn.pid
start(){
	/etc/init.d/myvpn enable
	/usr/sbin/myvpn start >/dev/null 2>&1 &
}

stop(){
	kill -9 `cat $PID` >/dev/null 2>&1
}

